# debian/aptitude-robot.cron.daily
#!/bin/sh
[ -f /etc/default/aptitude-robot ] && . /etc/default/aptitude-robot

if [ "$RUN_DAILY" = "no" ]; then
    exit
fi

SESSION_SCRIPT=/usr/share/aptitude-robot/aptitude-robot-session
[ -x $SESSION_SCRIPT ] && $SESSION_SCRIPT
